# LaunchDaemons/com.arcboxlabs.desktop.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arcboxlabs.desktop.daemon</string>

    <key>BundleProgram</key>
    <string>Contents/Helpers/com.arcboxlabs.desktop.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>com.arcboxlabs.desktop.daemon</string>
        <string>--docker-integration</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <!-- Kept for backward compatibility: older daemon versions write to
         stdout/stderr. Once the bundled daemon uses tracing-appender
         (arcbox >= 0.4), these can be removed. Paths moved from /tmp/
         to ~/.arcbox/log/ so all logs are in one place. launchd routes
         output to the system log if these keys are absent. -->
    <key>StandardOutPath</key>
    <string>/tmp/arcbox-daemon.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/arcbox-daemon.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
